# roles/polywiki_tomcat/templates/tomcat-wiki.service.j2
[Unit]
Description=Apache Tomcat instance for polywiki (wiki)
Documentation=https://tomcat.apache.org/
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
User={{ polywiki_tomcat_owner }}
Group={{ polywiki_tomcat_group }}
Environment="CATALINA_HOME={{ polywiki_tomcat_home }}"
Environment="CATALINA_BASE={{ polywiki_tomcat_base_dir }}"
Environment="CATALINA_PID={{ polywiki_tomcat_pid_file }}"
Environment="JAVA_HOME={{ polywiki_java_home }}"
Environment="CATALINA_OPTS=-Djava.net.preferIPv4Stack=true"
Environment="UMASK=0007"
WorkingDirectory={{ polywiki_tomcat_base_dir }}

# Utilise catalina.sh run (foreground) au lieu de start (forking)
# Plus propre sous systemd, les logs vont directement dans journald
ExecStart={{ polywiki_tomcat_home }}/bin/catalina.sh run

# Arrêt gracieux via le port de shutdown
ExecStop={{ polywiki_tomcat_home }}/bin/catalina.sh stop -force
TimeoutStopSec=30

# Si le processus crash, redémarrer après 10 secondes
Restart=on-failure
RestartSec=10

# Sécurité
NoNewPrivileges=true
PrivateTmp=true

# Limites de ressources (optionnel, ajuster selon besoins)
# LimitNOFILE=65536
# LimitNPROC=4096

[Install]
WantedBy=multi-user.target

